# com.zero.morning-brief.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Install with: cp com.zero.morning-brief.plist ~/Library/LaunchAgents/ && launchctl load ~/Library/LaunchAgents/com.zero.morning-brief.plist -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.zero.morning-brief</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/forge/getzero-os/.venv/bin/python3</string>
        <string>scanner/v6/card_push.py</string>
        <string>--morning</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/forge/getzero-os</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/zero-morning-brief.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/zero-morning-brief.err</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
